# recipes-core/udev/files/05-decrypt.rules
ACTION=="add", KERNEL=="mmcblk2p*", RUN+="/usr/local/bin/decryptdrive '%E{DEVNAME}'"
